# openrc/servicio-web
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
    need net
    need localmount
    need nfsmount
    need servicio-db-correo
    need servicio-db-blog
    after sshd
}

start() {
	ebegin "Arrancando servicios: Web, Wordpess"

        su luis -c "cd /Apps/docker/servicio-web && /usr/local/bin/fig up -d"

        eend $?
}

stop() {
	ebegin "Parando servicios: Web, Wordpress"

        su luis -c "cd /Apps/docker/servicio-web &&  /usr/local/bin/fig stop"

	eend $?
}

status() {

        su luis -c "cd /Apps/docker/servicio-web &&  /usr/local/bin/fig ps"

	eend $?
}
